# timetable-builder.cron
0 1 * * * root /opt/timetable-data-builder/timetableDataBuilder.sh >> /cronlogs/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job